# decrypted_and_extracted/dwr956_v1.0.0.7_r02_nb_p02003/root/usr/lib/opkg/info/base-files.list
/etc/cmv_batch
/etc/init.d/ipsec
/etc/rc.d/rc.bringup_lan
/etc/init.d/wlan_bringup_usb.sh
/etc/init.d/sshd_start
/etc/openwrt_version
/etc/tr69
/ramdisk_copy/etc/inetd.conf
/ramdisk_copy/tmp/u_lock
/etc/resolv.conf
/etc/rc.d/udhcpc_ad
/etc/init.d/ltq_switch_config.sh
/etc/rc.d/rc.bringup_lan_start
/var
/ramdisk_copy/var/lib/misc/udhcpd.leases
/etc/nsswitch.conf
/lib/firstboot/99_10_no_fo_cleanup
/etc/hotplug2-init.rules
/usr/local/bin/routing.sh
/etc/sys.conf.device
/etc/init.d/wifiStart
/etc/rc.d/ppa_config.sh
/etc/preinit
/etc/rc.d/rc.bringup_lan_ipv6_start
/lib/firstboot/50_pivot
/etc/rc.d/init.d
/etc/init.d/lanStart
/sbin/mount_root
/lib/firstboot/30_reset_copy_rom
/etc/host.conf
/etc/xl2tpd/xl2tpd.conf
/lib/upgrade/platform.sh
/etc/hotplug2-common.rules
/usr/sbin/patches_list.sh
/usr/sbin/version.sh
/etc/rc.d/get_wan_if
/etc/shells
/etc/init.d/cli
/lib/preinit/70_initramfs_test
/usr/local/bin/DBInit
/lib/preinit/50_choose_console
/lib/firstboot/10_no_fo_clear_overlay
/etc/protocols
/etc/init.d/appInit
/etc/init.d/mountd
/usr/local/bin/filter.sh
/etc/rc.d/ppa_up_config
/lib/preinit/20_check_jffs2_ready
/usr/local/bin/iptables_fw.sh
/lib/preinit/10_indicate_preinit
/ramdisk_copy/etc/udhcpd.conf
/ramdisk_copy/etc/TZ
/etc/init.d/load_gphy_firmware.sh
/etc/rc.d/ltq_port_mirror_stop.sh
/etc/hotplug.d/iface/00-netstate
/etc/ripd.conf
/usr/sbin/mini_fo.sh
/etc/dnrd
/etc/udhcpd.conf
/lib/firstboot/99_10_with_fo_cleanup
/usr/local/bin/lanServicesInit
/etc/rc.d/rc.bringup_lan_ipv6_stop
/sbin/ifdown
/etc/init.d/unload_ppa_modules.sh
/etc/switchports.conf
/usr/local/bin/firewall.sh
/etc/config/system
/usr/share/udhcpc/default.script
/etc/rc.d/S00jnr_init.sh
/etc/init.d/ltq_load_voice_tapi.sh
/usr/local/bin/udhcpc_usb0.script
/usr/local/bin/wan.dialup
/mnt/data/fs
/etc/init.d/inetd_start
/etc/hotplug.d/usb/ltq-usb-led.hotplug
/etc/init.d/led_config
/etc/ipsec.conf
/usr/sbin/bringup_single_mii_1.sh
/usr/local/bin/SessionLimit.sh
/etc/rc.d/rc.bringup_lan_ipv6
/etc/fstab
/lib/network/config.sh
/etc/init.d/pmcu.sh
/lib/preinit/30_device_fs_daemons
/etc/init.d/lte_daemon.sh
/lib/preinit/90_mount_no_jffs2
/etc/sysctl.conf
/sbin/sysupgrade
/lib/preinit/70_pivot_jffs2_root
/lib/preinit/99_10_run_init
/lib/firstboot/20_has_mini_fo
/lib/firstboot/10_reset_has_mini_fo
/etc/device_table.txt
/etc/init.d/demo_setup
/lib/firstboot/10_determine_parts
/sbin/led.sh
/etc/init.d/powered_up
/etc/rc.d/rc.bringup_lan_stop
/etc/hotplug.d/net/10-net
/etc/init_d5.script
/lib/preinit/50_indicate_regular_preinit
/lib/firstboot/40_no_fo_copy_ramoverlay
/etc/init.d/load_ppa_modules.sh
/usr/local/bin/104xfwdl
/etc/init.d/cronInit
/usr/local/bin/ap_fw_upgrade.sh
/etc/init.d/dropbear
/lib/preinit/10_essential_fs
/etc/rc.common
/usr/local/bin/QosRule.sh
/etc/init.d/ifx_load_ifxos_drv.sh
/ramdisk_copy/etc/snmp/snmpd.conf
/etc/init.d/syslog-ng
/etc/init.d/samba
/lib/preinit/10_indicate_failsafe
/lib/preinit/40_init_shm
/etc/rc.d/wlan_discover
/etc/init.d/systemStart
/init
/usr/sbin/bringup_single_mii_2.sh
/etc/init.d/sshd_stop
/ramdisk_copy/etc/dsl_api/cmv_scripts
/etc/spmd.pwd
/etc/init.d/passwd.sh
/etc/ipsec.secrets
/usr/local/bin/update_to_default_config_essid.sh
/etc/rc.d/get_lan_if
/ramdisk_copy/etc/ppp/resolv.conf
/etc/init.d/inetd
/etc/modules.conf
/etc/init.d/lan_services_bringup.sh
/etc/init.d/smp_irq_affinity.sh
/etc/init.d/ltq_show_version.sh
/lib/preinit/40_mount_jffs2
/usr/local/bin/udhcpc_wan.script
/lib/preinit/20_device_fs_mount
/lib/preinit/80_mount_root
/etc/init.d/wlan_bringup.sh
/etc/rc.local
/usr/local/bin/firewall_rule.sh
/etc/rc.d/ltq_switch_functions.sh
/etc/rc.d/rebootcpe.sh
/etc/rc.d/invoke_upgrade.sh
/usr/local/bin/dmz.sh
/ramdisk_copy/flash/print_mcast_debug.sh
/lib/firstboot/30_is_rootfs_mounted
/usr/local/bin/upnp_service.sh
/sbin/hotplug-call
/etc/snmp
/etc/init.d/wanStart
/etc/fw_env.config
/lib/modules/2.6.32.42/GPL_NetUSB.ko
/etc/rc.d/backup
/etc/crontabs/root
/firmware
/etc/rc.d/config.sh
/ramdisk_copy/etc/ilmid/ilmid.conf
/etc/rc.d/rc.bringup_services_stop
/etc/rc.d/ltq_port_mirror_start.sh
/etc/rc.d/rc.bringup_services_start
/etc/init.d/sshd
/etc/init.d/iptables
/usr/sbin/factorycfg.sh
/lib/preinit/40_mount_devpts
/sbin/modprobe.sh
/lib/preinit/90_init_console
/etc/rc.d/killproc
/usr/local/bin/dhcpd_br1Init
/usr/local/bin/forward.sh
/etc/services
/lib/preinit/99_10_mount_no_mtd
/etc/init_a5.script
/lib/preinit/10_check_for_mtd
/usr/local/bin/lanInit
/etc/rc.d/checkpid
/etc/init.d/ltq_load_voice_vmmc.sh
/lib/upgrade/common.sh
/etc/rc.d/ltq_switch_conf_examples.sh
/etc/rc.d/restore_tr69_defaults
/etc/rc.d/ltq_vlan_sep.sh
/usr/local/bin/QosInit.sh
/etc/sysupgrade.conf
/lib/firstboot/05_firstboot_skip
/etc/rc.d/rc.bringup_ipv6_staticRoutes
/etc/ppp/options.pptpd
/usr/local/bin/urlfilter.sh
/etc/init.d/ftpd
/flash
/etc/init.d/ifx_load_voip_common_drivers.sh
/lib/upgrade/keep.d/base-files
/ramdisk_copy/etc/services
/ramdisk_copy/flash/env.sh
/etc/ppp/options.pptp
/etc/init.d/version.sh
/lib/preinit/40_run_failsafe_hook
/etc/group
/etc/rc.d/rc.getBridgeStats
/lib/preinit/30_failsafe_wait
/etc/init.d/httpd
/etc/rc.d/rc.bringup_voip_start
/etc/init.d/PPA.sh
/etc/functions.sh
/etc/inetd.conf
/sbin/firstboot
/ramdisk_copy/flash/print_debug_info.sh
/etc/init.d/dwatch
/etc/rc.d/udhcpc.script
/etc/pptpd.conf
/ramdisk_copy/flash/lock.txt
/bin/login.sh
/etc/rc.d/rc.bringup_wlan
/etc/init.d/ipqos_qprio_cfg.sh
/ramdisk_copy/tmp/s_lock
/ramdisk_copy/etc/passwd
/etc/zebra.conf
/etc/ppp/options
/usr/local/bin/QosPolicy.sh
/etc/hosts
/etc/hotplug.d/usb/10-usb
/usr/local/bin/vpn_method.sh
/lib/preinit/60_init_hotplug
/etc/init.d/ltq_load_voice_kpi2udp.sh
/etc/rc.d/rc.restart_voip.sh
/usr/local/bin/fwnat.sh
/etc/init.d/minissdpd.sh
/etc/config/network
/usr/local/bin/ddns.sh
/etc/rc.d/rc.bringup_staticRoutes
/etc/init.d/rcS
/usr/lib/common.awk
/usr/sbin/savecfg.sh
/sbin/ifup
/etc/init.d/telnetd
/rom/note
/etc/inittab
/etc/openwrt_release
/lib/preinit/05_mount_skip
/etc/rc.d/rc.bringup_services
/sbin/wifi
/etc/ppp/chap-secrets
/etc/rc.conf
/usr/local/bin/schedule_rule.sh
/etc/ilmid
/etc/profile
/usr/local/bin/led_config.sh
/etc/init.d/disablePHY
/etc/hotplug.d/iface/10-routes
/etc/init_d5_2.script
/usr/local/bin/remmag.sh
/lib/modules/2.6.32.42/NetUSB.ko
/lib/preinit/99_10_failsafe_login
/etc/TZ
/etc/init.d/wlan_load.sh
/etc/rc.d/rc.bringup_voip_stop
/lib/firstboot/20_no_fo_mount_jffs
/usr/local/bin/dhcpdInit
/etc/init.d/tftpd
/etc/ppp/options.xl2tpd
/etc/rc.d/free_memory.sh
/lib/preinit/90_restore_config
/lib/upgrade/keep.d/base-files-essential
/lib/firstboot/30_no_fo_pivot
/etc/init.d/ethwanStart
/etc/init.d/enablePHY
/etc/passwd
/sbin/hotplug
/etc/sys.conf
/etc/hotplug.d/ieee1394/10-ieee1394
/etc/init.d/systemFinish
/etc/init_a5_2.script
/usr/sbin/vol_mgmt
/etc/mtab
/usr/local/etc/VID_PID_LIST.dat
/lib/firstboot/40_copy_ramoverlay
/etc/init.d/boot
/ramdisk_copy/tmp/system_status
/etc/init.d/usb_host_device.sh
/etc/rc.d/get_wlan_stats
/ramdisk_copy/etc/ripd.conf
/etc/ppp/peers
/lib/functions/boot.sh
/ramdisk_copy/etc/hosts
/etc/rc.d/rc.getBridgeIfs
/etc/ppp/resolv.conf
/bin/ipcalc.sh
/etc/init.d/init_ipqos.sh
/mnt/usb
/etc/diag.sh
/etc/hostapd.conf
/etc/rc.d/udhcpc_lan.script
/etc/init.d/udhcpc
/etc/rc.d/rc.getIfStats
/lib/firstboot/20_reset_clear_jffs
/etc/codeclist
/usr/local/bin/ntp.sh
/lib/preinit/41_merge_overlay_hooks
