# decrypted_and_extracted/dwr956_v1.0.0.7_r02_nb_p02003/root/usr/lib/opkg/info/ifx-voice-cpe-tapi.list
/lib/modules/2.6.32.42/drv_tapi.ko
/etc/init.d/ltq_load_voice_tapi.sh
/opt/lantiq/bin/drv_tapi.ko
