# decrypted_and_extracted/dwr956_v1.0.0.9_r06_nb_p02003/root/usr/lib/opkg/info/dnsmasq.list
/etc/dnsmasq.conf
/etc/config/dhcp
/etc/init.d/dnsmasq
/usr/sbin/dnsmasq
