# decrypted_and_extracted/dwr956_v1.0.0.9_r06_nb_p02003/root/usr/lib/opkg/info/dropbear.list
/usr/bin/dbclient
/usr/bin/scp
/etc/init.d/dropbear
/usr/bin/dropbearkey
/etc/config/dropbear
/usr/sbin/dropbear
/lib/upgrade/keep.d/dropbear
/usr/bin/ssh
/etc/dropbear
