# decrypted_and_extracted/dwr956_v1.0.0.9_r06_nb_p02003/root/usr/lib/opkg/info/ifx-voice-cpe-vmmc.list
/lib/modules/2.6.32.42/drv_vmmc.ko
/etc/init.d/ltq_load_voice_vmmc.sh
/opt/lantiq/bin/drv_vmmc.ko
