# decrypted_and_extracted/dwr956_v1.0.0.9_r06_nb_p02003/root/usr/lib/opkg/info/ifx-voip-common.list
/lib/modules/2.6.32.42/voip_timer_driver.ko
/etc/init.d/ifx_load_voip_common_drivers.sh
