# decrypted_and_extracted/dwr956_v1.0.0.9_r06_nb_p02003/root/usr/lib/opkg/info/ltq-dwatch.list
/usr/sbin/dwatchctl
/usr/sbin/dwatch
/etc/init.d/dwatch
