# decrypted_and_extracted/dwr956_v1.0.0.9_r06_nb_p02003/root/usr/lib/opkg/info/mountd.list
/etc/init.d/mountd
/etc/init.d/usb_automount_status.sh
/etc/config/mountd
/sbin/mountd
/etc/init.d/usb_led.sh
