# decrypted_and_extracted/dwr956_v1.0.0.7_r02_nb_p02003/root/etc/init.d/acinit.sh
#!/bin/sh /etc/rc.common

START=992

start() 
{
#temp sol : put the wifi AC intial script here
    [ -e /usr/local/bin/wifiACInit ]  && /usr/local/bin/wifiACInit
}



