# decrypted_and_extracted/dwr956_v1.0.0.7_r02_nb_p02003/root/etc/init.d/boot
#!/bin/sh /etc/rc.common

START=20

start() {
	udevd --daemon
	load_modules /etc/modules.d/*
}
