# decrypted_and_extracted/dwr956_v1.0.0.7_r02_nb_p02003/root/etc/init.d/led_config
#!/bin/sh /etc/rc.common

START=68

start() {
	/usr/local/bin/led_config.sh &
}

