# decrypted_and_extracted/dwr956_v1.0.0.8/root/etc/init.d/disablePHY
#!/bin/sh /etc/rc.common

START=11

start() {
	echo "disable PHY..."
	switch_utility PortLinkCfgSet 0 0 0 0 10 1 1 0 0 0
	switch_utility PortLinkCfgSet 1 0 0 0 10 1 1 0 0 0
	switch_utility PortLinkCfgSet 2 0 0 0 10 1 1 0 0 0
	switch_utility PortLinkCfgSet 4 0 0 0 10 1 1 0 0 0
}
