# decrypted_and_extracted/dwr956_v1.0.0.8/root/etc/init.d/lanStart
#!/bin/sh /etc/rc.common

START=65

start() {
	/etc/init.d/dwatch start
	[ -e /usr/local/bin/lanInit ]                && lanInit
}
