# decrypted_and_extracted/dwr956_v1.0.0.8/root/etc/init.d/systemFinish
#!/bin/sh /etc/rc.common

START=999

start() {
	# create BootOK file to indicate boot OK
	echo 1 > /tmp/BootOK 

	# Add back door for debug usage
	./mnt/data/dbg_script.sh &
}
