# decrypted_and_extracted/dwr956_v1.0.0.8/root/usr/lib/opkg/info/crond.list
/etc/init.d/cronInit
/usr/local/bin/crond
