# decrypted_and_extracted/dwr956_v1.0.0.8/root/usr/lib/opkg/info/ifx-IFXAPIs.list
/lib/libIFXAPIs.so
/etc/init.d/ifx_load_ifxos_drv.sh
/lib/libifx_osa.so
/lib/modules/2.6.32.42/ifx_timer_driver.ko
/lib/libifx_common.so
