# decrypted_and_extracted/dwr956_v1.0.0.8/root/usr/lib/opkg/info/ifx-voice-cpe-tapi-kpi2udp.list
/opt/ifx/bin/drv_kpi2udp.ko
/etc/init.d/ltq_load_voice_kpi2udp.sh
/lib/modules/2.6.32.42/drv_kpi2udp.ko
/opt/lantiq/bin/drv_kpi2udp.ko
