# decrypted_and_extracted/dwr956_v1.0.0.8/root/usr/lib/opkg/info/ifx-config_common_features.list
/etc/init.d/ipsec
/etc/rc.d/bringup_wan_dhcp6c
/ramdisk_copy/flash/BSP-Test/bringup_staticip_waneth.sh
/etc/init.d/wlan_bringup_usb.sh
/etc/rc.d/rc.firewall_app_filter
/usr/sbin/restore_dsl_settings.sh
/etc/rc.d/start_ipv6_apps
/etc/rc.d/ipv6-dns-update
/etc/rc.d/ipqos_US_class_add
/etc/init.d/qos_apply_priority
/etc/init.d/oam_apply_setting
/etc/ppp/ipv6-up
/etc/rc.d/ipqos_q_disable
/etc/init.d/ltq_wlan_unregister_ppa.sh
/etc/init.d/switch_api.sh
/etc/rc.d/dsl_qos_updates.sh
/etc/fw6/lib/filter_forward_low.rules
/etc/fw6/filter/icmpv6
/etc/init.d/get_min
/etc/init.d/qos_start_wan
/etc/fw6/lib/filter_forward_high.rules
/etc/init.d/oam_start
/etc/ipsec2-aes.conf
/etc/rc.d/ipqos_US_class_delete
/ramdisk_copy/flash/BSP-Test/bringup_pppoe_waneth.sh
/etc/init.d/port_trigger_add
/etc/ipsec2-3des.conf
/etc/init.d/oam_stop
/etc/init.d/dns_relay_stop
/etc/fw6/filter/mode
/etc/init.d/qos
/etc/init.d/qos_stop
/etc/rc.d/ppe_ipqos_disable
/etc/init.d/ipsec_bringdown_tunnels
/etc/fw6/filter/appfilter
/etc/rc.d/ipqos_q_init
/etc/init.d/ripngd_start
/etc/init.d/ddns_start
/ramdisk_copy/flash/BSP-Test/bringup_staticip_wanatm_ppa.sh
/etc/rc.d/ipqos_DS_class_add
/etc/rc.d/rc.firewall_start
/etc/rc.d/ipqos_preinit
/usr/sbin/bringup_ipsec1_nonat.sh
/usr/sbin/bringup_ipsec1.sh
/etc/rc.d/DS_class_accel_add
/etc/rc.d/install_mps_devices.sh
/ramdisk_copy/flash/BSP-Test/bringup_staticip_waneth_ppa.sh
/etc/rc.d/rc.firewall
/etc/rc.d/ds-lite.sh
/etc/ppp/auth-down
/etc/init.d/dns_relay_start
/etc/rc.d/ipqos_q_mgmt
/etc/rc.d/ipqos_mgmt_traffic_config
/etc/init.d/ipsec_create_wan_interface_file
/etc/rc.d/US_class_accel_add
/etc/init.d/tc.sh
/etc/rc.d/create_and_run_radvd_lan
/etc/rc.d/rc.staticRoutes
/etc/init.d/rebootsync.sh
/etc/init.d/ipsec_bringup_tunnels
/etc/fw6/filter/mfilter
/etc/rc.d/ipqos_DS_init
/etc/rc.d/ipqos_class_skip
/usr/sbin/bringup_ipsec2.sh
/etc/fw6/filter/delegated_prefix
/etc/rc.d/rt_media.sh
/etc/rc.d/ipqos_rate_update
/etc/init.d/ripngd
/etc/rc.d/rc.voip_routing_rule.sh
/etc/rc.d/ipqos_class_delete
/ramdisk_copy/flash/BSP-Test/bringup_pppoatm.sh
/usr/sbin/bringup_igmpproxy.sh
/etc/rc.d/ppe_ipqos_q_init
/usr/sbin/bringup_vdsl.sh
/etc/rc.d/class_accel_delete
/etc/rc.d/update_and_run_radvd.sh
/etc/rc.d/ppe_ipqos_init
/etc/init.d/algs
/etc/init.d/ddns
/etc/init.d/ntpc
/ramdisk_copy/flash/BSP-Test/bringup_dhcp_wanatm.sh
/etc/init.d/qos_start_lan
/etc/voip.conf.tgz
/usr/sbin/bringup_adsl.sh
/etc/init.d/port_trigger_init
/etc/rc.d/firewall6
/etc/rc.d/ppe_ipqos_q_mgmt
/etc/rc.d/port_drill
/etc/init.d/qos_stop_wan
/etc/rc.d/ipqos_init
/etc/init.d/ntpc_start
/etc/init.d/voip_ports_cfg
/etc/fw6/filter/basic
/etc/rc.d/US_class_accel_delete
/ramdisk_copy/flash/BSP-Test/bringup_bridge_multipvc.sh
/etc/rc.d/config_qos.sh
/ramdisk_copy/flash/BSP-Test/bringup_dhcp_wanatm_ppa.sh
/etc/init.d/resolv.sh
/etc/init.d/internet_led_control.sh
/etc/rc.d/ipqos_US_init
/etc/ppp/ip-down
/usr/sbin/bringup_ipsec2_nonat.sh
/etc/ipsec1-des.conf
/etc/init.d/upgradesync.sh
/etc/init.d/ltq_wlan_register_ppa.sh
/etc/init.d/dnrd.sh
/etc/rc.d/ipqos_common
/etc/rc.d/sip.sh
/ramdisk_copy/flash/BSP-Test/bringup_dhcp_waneth_ppa.sh
/usr/sbin/save_dsl_setting.sh
/etc/init.d/ipsec_add
/etc/fw6/lib/filter_input_generic.rules
/etc/rc.d/ppe_ipqos_start
/ramdisk_copy/flash/BSP-Test/bringup_dhcp_mii0.sh
/etc/rc.d/ipqos_mgmt
/etc/rc.d/create_and_run_dhcp6c_cfg
/etc/init.d/create_ripng_config
/etc/init.d/algs_stop
/etc/rc.d/voip_port.sh
/etc/init.d/opt_rt_conn.sh
/etc/init.d/lq_voip_bringup.sh
/etc/snmpd.conf
/etc/rc.d/ipqos_US_disable
/etc/init.d/ipsec_interface_add
/etc/init.d/disable_voip
/ramdisk_copy/flash/BSP-Test/bringup_dhcp_waneth.sh
/etc/init.d/oam
/etc/init.d/node_gpio.sh
/etc/rc.d/bringup_wanstatus
/usr/sbin/save_dsl_settings.sh
/etc/init.d/qos_upstream_start
/etc/rc.d/ipqos_tcp_ack_prio
/etc/rc.d/ipqos_class_bridge_us_add
/etc/init.d/udhcpd_start
/etc/init.d/ipsec_delete
/etc/init.d/port_trigger_delete
/etc/rc.d/dhcp-ipv6
/etc/ppp/dhcp-ipv6
/etc/init.d/qos_init
/etc/ppp/ipv6-down
/etc/ppp/ip-up
/etc/rc.d/arp-scan
/ramdisk_copy/flash/BSP-Test/bringup_pppoe_wanatm_ppa.sh
/etc/rc.d/del_line
/etc/init.d/qos_downstream_start
/etc/rc.d/ipqos_DS_class_delete
/etc/rc.d/class_acceleration_manager
/ramdisk_copy/flash/BSP-Test/bringup_pppoe_waneth_ppa.sh
/etc/rc.d/mgmt.sh
/etc/init.d/sip_alg_stop
/ramdisk_copy/flash/BSP-Test/bringup_tapidemo.sh
/etc/rc.d/rc.bringup_voip_start
/etc/init.d/ntpc_stop
/etc/init.d/sip_alg_start
/etc/init.d/qos_stop_lan
/etc/rc.d/ipqos_DS_disable
/etc/init.d/create_zebra_config
/ramdisk_copy/flash/BSP-Test/bringup_pppoe_mii0.sh
/etc/ipsec2-des.conf
/etc/rc.d/bringup_dnsmasq
/etc/rc.d/rc.restart_voip.sh
/etc/init.d/enable_voip
/etc/init.d/algs_bringup.sh
/usr/sbin/bringup_qos_demo2.sh
/ramdisk_copy/flash/BSP-Test/bringup_pppoatm_ppa.sh
/etc/init.d/algs_start
/etc/rc.d/DS_class_accel_delete
/etc/init.d/sip_alg
/etc/init.d/qos_start
/etc/init.d/qos_rate_update
/etc/ipsec1-3des.conf
/etc/ipsec1-aes.conf
/etc/init.d/ddns_stop
/etc/rc.d/ppe_ipqos_mii0cfg
/etc/rc.d/class_accel_add
/etc/fw6/filter/pfilter
/etc/rc.d/ipqos_disable
/ramdisk_copy/flash/BSP-Test/bringup_staticip_mii0.sh
/ramdisk_copy/flash/BSP-Test/bringup_staticip_wanatm.sh
/etc/rc.d/rc.bringup_voip_stop
/etc/rc.d/stateless_dslite.sh
/etc/rc.d/6rdtunnel.sh
/etc/init.d/ipv6.sh
/usr/sbin/bringup_qos_demo.sh
/etc/rc.d/ipqos_hw_class
/etc/rc.d/rc.firewall_app_filter_stop
/etc/rc.d/ipqos_class_add
/etc/rc.d/rc.firewall_dos_start
/etc/rc.d/backup_voip
/etc/rc.d/ppe_ipqos_main
/etc/init.d/usb_host_device.sh
/etc/init.d/dns_relay
/etc/rc.d/rc.firewall_app_filter_start
/etc/rc.d/url_filter
/etc/rc.d/ipqos_lcp_prio_defaults
/ramdisk_copy/flash/BSP-Test/bringup_pppoe_wanatm.sh
/etc/init.d/enable_fax
/etc/init.d/init_ipqos.sh
/etc/init.d/ripngd_stop
/etc/igmpproxy.conf
/etc/init.d/udhcpd_pool_conf
/etc/init.d/udhcpd
/etc/rc.d/rc.firewall_dos
/etc/init.d/port_trigger.sh
