# decrypted_and_extracted/dwr956_v1.0.0.8/root/usr/lib/opkg/info/syslog-ng.list
/etc/init.d/syslog-ng
/etc/syslog-ng/syslog-ng.conf
/usr/sbin/syslog-ng
