# decrypted_and_extracted/DWR956_v1.0.0.11_r03_NB_P02003/root/usr/lib/opkg/info/atheros.list
/etc/rc.d/ath_wlan_get_ap_dyn_info
/etc/rc.d/ath_wlan_wps_get_pin
/root/.dummy
/root/.profile
/etc/rc.d/ath_wlan_mac_ctrl_modify
/etc/rc.d/ath_wlan_get_stats
/etc/ath/WSC_sta.conf
/sbin/pktlogconf
/etc/ath/RADIUS_MAC.ap_bss
/usr/local/bin/wifiACACL.sh
/etc/init.d/acinit.sh
/etc/ath/apcfg
/etc/rc.d/ath_wlan_wps_generate_pin
/etc/rc.d/ath_wlan_remove_vap
/root/.configData
/lib/modules/2.6.32.42/net/umac.ko
/etc/rc.d/ath_wlan_uninit
/sbin/wpscli_sta
/lib/modules/2.6.32.42/net/ath_pktlog.ko
/etc/rc.d/ath_wlan_wmm_modify
/etc/rc.d/rc.bridge
/etc/ath/PSK.ap_bss
/etc/wpa2/.dummy
/sbin/wpscli_ap
/sbin/athstats
/etc/rc.d/ath_wlan_capability
/lib/libiw.so.29
/etc/ath/PSK.sta
/etc/rc.d/ath_wlan_ap_stop
/etc/rc.d/ath_wlan_wps_reset
/etc/rc.d/ath_wlan_get_sec_cntrls
/etc/rc.d/ath_wlan_ap_start
/etc/rc.d/ath_wlan_wps_trigger_pin_pairing
/etc/ath/killVAP
/etc/ath/apup
/sbin/hostapd_cli_ath
/sbin/hostapd_ath
/lib/modules/2.6.32.42/net/ath_dfs.ko
/sbin/pktlogdump
/sbin/80211stats
/etc/ath/activateVAP
/etc/rc.d/ath_wlan_wps_get_profile
/etc/rc.d/ath_wlan_get_wps_dyn_info
/etc/rc.d/ath_wlan_wps_config
/etc/rc.d/ath_wlan_get_assoc
/sbin/cfg
/etc/rc.d/ath_wlan_get_radio_dyn_info
/lib/modules/2.6.32.42/net/ath_dev.ko
/etc/rc.d/rc.wlan
/lib/modules/2.6.32.42/net/asf.ko
/etc/rc.d/ath_wlan_wps_trigger_pbc_pairing
/sbin/athstatsclr
/etc/rc.d/ath_wlan_wps_reset_pin
/usr/local/bin/wifiACInit
/etc/ath/EAP.ap_bss
/etc/ath/makeVAP
/etc/ath/apdown
/lib/modules/2.6.32.42/net/adf.ko
/etc/ath/WSC.conf
/etc/rc.d/ath_wlan_init
/lib/modules/2.6.32.42/net/ath_hal.ko
/usr/local/bin/wifiACHwPBC.sh
/usr/local/bin/wifiACWPS.sh
/etc/rc.d/rc.network
/sbin/radartool
/etc/rc.d/ath_wlan_get_wps_regs_dyn_info
/etc/rc.d/ath_wlan_sec_modify
/lib/modules/2.6.32.42/net/ath_rate_atheros.ko
/sbin/wlanconfig
/etc/rc.d/ath_wlan_wps_get_status
