# decrypted_and_extracted/DWR956_v1.0.0.11_r03_NB_P02003/root/usr/lib/opkg/info/ifx-dect.list
/usr/sbin/DECT_VERSION
/usr/drivers/drv_dect.ko
/opt/ifx/downloads/BMCFw.BIN
/usr/drivers/paging.ko
/opt/ifx/downloads/COSICFw.BIN
/usr/drivers/drv_timer.ko
/etc/init.d/ifx_load_dect_drivers.sh
/usr/sbin/DectCli
