# decrypted_and_extracted/DWR956_v1.0.0.11_r03_NB_P02003/root/usr/lib/opkg/info/ltq-base-files-common.list
/etc/banner
/etc/init.d/ltq_show_version.sh
/opt/lantiq/bin/inst_driver.sh
/opt/lantiq/bin/whatz.sh
/opt/lantiq/bin/what.sh
/opt/lantiq/bin/debug_level.cfg
