# decrypted_and_extracted/DWR956_v1.0.0.11_r03_NB_P02003/root/usr/lib/opkg/info/openswan.list
/etc/init.d/ipsec
/usr/lib/ipsec/_startklips
/usr/libexec/ipsec/tncfg
/usr/libexec/ipsec/_pluto_adns
/usr/lib/ipsec/_updown.mast
/usr/libexec/ipsec/rsasigkey
/usr/libexec/ipsec/secrets
/lib/upgrade/keep.d/ipsec
/usr/libexec/ipsec/spigrp
/usr/libexec/ipsec/barf
/usr/libexec/ipsec/showdefaults
/usr/lib/ipsec/_secretcensor
/usr/libexec/ipsec/ikeping
/usr/libexec/ipsec/newhostkey
/usr/libexec/ipsec/pf_key
/usr/lib/ipsec/_updown.klips
/usr/lib/ipsec/_include
/usr/libexec/ipsec/ranbits
/usr/libexec/ipsec/auto
/usr/libexec/ipsec/policy
/usr/libexec/ipsec/showhostkey
/usr/lib/ipsec/_realsetup
/usr/lib/ipsec/_plutorun
/usr/libexec/ipsec/whack
/usr/libexec/ipsec/setup
/usr/libexec/ipsec/spi
/usr/lib/ipsec/_updown
/usr/libexec/ipsec/eroute
/usr/libexec/ipsec/addconn
/lib/upgrade/keep.d/openswan
/usr/lib/ipsec/_startnetkey
/usr/lib/ipsec/_keycensor
/usr/lib/ipsec/_plutoload
/usr/libexec/ipsec/verify
/usr/sbin/ipsec
/usr/lib/ipsec/_copyright
/usr/libexec/ipsec/klipsdebug
/usr/libexec/ipsec/pluto
/usr/libexec/ipsec/look
/usr/lib/ipsec/_updown.netkey
