# decrypted_and_extracted/DWR956_v1.0.0.11_r03_NB_P02003/root/usr/lib/opkg/info/quagga.list
/usr/sbin/quagga.init
/usr/sbin/zebra
/etc/quagga/zebra.conf
/usr/sbin/watchquagga
/etc/init.d/quagga
